# crs_automate_stop.service
[Unit]
Description=Stop crs_automate.service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl stop crs_automate.service